# omnibus/package-scripts/nodex-agent-deb/postrm
#!/bin/sh
# Perform necessary cleanup steps after nodex-agent package is removed, purged, upgraded, or deconfigured.

INSTALL_DIR=/opt/nodex-agent
SERVICE_NAME=nodex-agent
HOME_DIR=/home/nodex/

set -e

case "$1" in
    remove)
        echo "Stopping service $SERVICE_NAME..."
        if command -v systemctl >/dev/null 2>&1; then
            if ! systemctl stop "$SERVICE_NAME"; then
                echo "Failed to stop service $SERVICE_NAME"
            fi
        fi

        echo "Disabling service $SERVICE_NAME..."
        if command -v systemctl >/dev/null 2>&1; then
            if ! systemctl disable "$SERVICE_NAME"; then
                echo "Failed to disable service $SERVICE_NAME"
            fi
            if ! systemctl daemon-reload; then
                echo "Failed to reload systemctl daemon"
            fi
        fi

        if [ -d "$INSTALL_DIR" ]; then
            echo "Removing installation directory $INSTALL_DIR..."
            if ! rm -rf "$INSTALL_DIR"; then
                echo "Failed to remove installation directory"
            fi
        fi
        ;;

    purge)
        echo "Purging configuration and removing directories..."

        # Remove the service, if it still exists
        if command -v systemctl >/dev/null 2>&1; then
            if ! systemctl stop "$SERVICE_NAME"; then
                echo "Failed to stop service $SERVICE_NAME"
            fi
            if ! systemctl disable "$SERVICE_NAME"; then
                echo "Failed to disable service $SERVICE_NAME"
            fi
            if ! systemctl daemon-reload; then
                echo "Failed to reload systemctl daemon"
            fi
        fi

        # Remove user and group
        if id -u nodex >/dev/null 2>&1; then
            echo "Removing nodex user and group..."
            if ! userdel -r nodex; then
                echo "Failed to remove nodex user"
            fi
            if ! groupdel nodex; then
                echo "Failed to remove nodex group"
            fi
        fi

        # Remove the installation directory
        if [ -d "$INSTALL_DIR" ]; then
            echo "Removing installation directory $INSTALL_DIR..."
            if ! rm -rf "$INSTALL_DIR"; then
                echo "Failed to remove installation directory"
            fi
        fi

        # Remove configuration files
        if [ -d "$HOME_DIR" ]; then
            echo "Removing configuration directory $HOME_DIR..."
            if ! rm -rf "$HOME_DIR"; then
                echo "Failed to remove configuration directory"
            fi
        fi
        ;;

    *)
esac

echo "$SERVICE_NAME post-removal script completed."
exit 0
