# omnibus/config/templates/init-scripts/systemd.service.erb
[Unit]
Description=Nodex Agent Service
After=network.target

[Service]
ExecStart=<%= install_dir %>/bin/nodex-agent
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
User=nodex
Group=nodex
PIDFile=<%= install_dir %>/var/run/nodex-agent.pid

PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
ReadOnlyPaths=<%= install_dir %>/bin

[Install]
WantedBy=multi-user.target
