# omnibus/config/templates/init-scripts/systemd.socket.erb
[Unit]
Description=Nodex Agent Socket

[Socket]
ListenStream=/home/nodex/run/nodex-agent.sock
SocketMode=0666

[Install]
WantedBy=sockets.target
