# omnibus/config/templates/init-scripts/systemd.socket.erb
[Unit]
Description=Nodex Agent Socket

[Socket]
ListenStream=/home/nodex/.nodex/run/nodex.sock
# Overridden by the umask setting, but just to be clear
SocketMode=0766
SocketUser=nodex
SocketGroup=nodex

[Install]
WantedBy=sockets.target
