# omnibus/config/templates/init-scripts/systemd.service.erb
[Unit]
# Service description and dependency management.
Description=Nodex Agent Service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple

# Core service execution settings.
ExecStart=/home/nodex/bin/nodex-agent

# Restart behavior: restarts on failure, with a 10-second delay between attempts.
# Limits restarts to 5 attempts within a 300-second period to avoid excessive restarts.
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Runs the service as the 'nodex' user and group, and manages the process ID.
User=nodex
Group=nodex

# Security and isolation settings:
# - PrivateTmp: Isolates temporary files from other services.
# - NoNewPrivileges: Prevents gaining additional privileges.
# - ProtectSystem, ProtectKernelModules, ProtectKernelTunables, ProtectControlGroups: Limits the service’s ability to modify the system, enhancing security.
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=full
ReadWritePaths=%h/.config
ReadWritePaths=%h/.nodex
ReadWritePaths=%h/bin
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true

[Install]
# Makes the service start automatically in multi-user mode.
WantedBy=multi-user.target
