# omnibus/config/templates/init-scripts/systemd.socket.erb
[Unit]
Description=Nodex Agent Socket

[Socket]
ListenStream=/home/nodex/run/nodex.sock
SocketMode=0660

[Install]
WantedBy=sockets.target
